# systemd/crypto-trading-be.service
[Unit]
Description=Crypto Trading Connector Backend
Documentation=https://github.com/Kohei-Sato-1221/crypto-trading-connector-by-kiro
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/workspace/crypto-trading-connector-by-kiro/crypto-trading-connector-be
Environment=PATH=/usr/local/go/bin:/usr/local/bin:/usr/bin:/bin
Environment=GOPATH=/home/pi/go
Environment=GOCACHE=/home/pi/.cache/go-build
Environment=DB_HOST=host
Environment=DB_USER=user
Environment=DB_PASSWORD=user
Environment=DB_NAME=db_name
Environment=DB_PORT=3307
Environment=SERVER_PORT=8090
Environment=BITFLYER_API_URL=https://api.bitflyer.com
Environment=BITFLYER_API_KEY=api_key
Environment=BITFLYER_API_SECRET=api_secret
Environment=CORS_ALLOWED_ORIGINS=http://10.221.1.162:3000,http://10.221.1.199:3000,http://localhost:3000
ExecStartPre=/usr/local/go/bin/go build -o bin/server cmd/server/main.go
ExecStart=/home/pi/workspace/crypto-trading-connector-by-kiro/crypto-trading-connector-be/bin/server
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crypto-trading-be

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/workspace/crypto-trading-connector-by-kiro
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target