# systemd/crypto-trading-be.service
[Unit]
Description=Crypto Trading Connector Backend
Documentation=https://github.com/Kohei-Sato-1221/crypto-trading-connector-by-kiro
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
# WorkingDirectory=/home/pi/workspace/crypto-trading-connector-by-kiro/crypto-trading-connector-be
EnvironmentFile=/home/pi/workspace/crypto-trading-connector-by-kiro/crypto-trading-connector-be/.env
ExecStart=/bin/bash -c 'cd /home/pi/workspace/crypto-trading-connector-by-kiro/crypto-trading-connector-be && ./server'
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crypto-trading-be

# Security settings (disabled for troubleshooting)
# NoNewPrivileges=false
# PrivateTmp=false
# ProtectSystem=false
# ProtectHome=false
UMask=0022

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target