# systemd/crypto-trading-fe.service
[Unit]
Description=Crypto Trading Connector Frontend
Documentation=https://github.com/Kohei-Sato-1221/crypto-trading-connector-by-kiro
After=network.target crypto-trading-be.service
Wants=network.target
Requires=crypto-trading-be.service

[Service]
Type=simple
User=pi
Group=pi
# WorkingDirectory=/home/pi/workspace/crypto-trading-connector-by-kiro/crypto-trading-connector-fe
Environment=PATH=/opt/nodejs/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/home/pi/workspace/crypto-trading-connector-by-kiro/crypto-trading-connector-fe/.env
Environment=NODE_ENV=production
ExecStartPre=/bin/bash -c 'cd /home/pi/workspace/crypto-trading-connector-by-kiro/crypto-trading-connector-fe && /opt/nodejs/bin/npm ci --only=production'
ExecStartPre=/bin/bash -c 'cd /home/pi/workspace/crypto-trading-connector-by-kiro/crypto-trading-connector-fe && /opt/nodejs/bin/npm run build'
ExecStart=/bin/bash -c 'cd /home/pi/workspace/crypto-trading-connector-by-kiro/crypto-trading-connector-fe && /opt/nodejs/bin/npm run preview -- --host 0.0.0.0 --port 3000'
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crypto-trading-fe

# Security settings (disabled for troubleshooting)
# NoNewPrivileges=false
# PrivateTmp=false
# ProtectSystem=false
# ProtectHome=false
UMask=0022

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target